# root/etc/services.d/app/run
#!/usr/bin/with-contenv bash

cd "${HOME}"

if [ "${USE_CELERY}" = "true" ]; then
    n=0
    m=5
    while true; do
        sleep $(((n**2+1)*5))
        s6-setuidgid abc \
            celery status --timeout=5 >/dev/null 2>&1 && break
        ((n++))
        if [ "$n" -gt "$m" ]; then
            s6-svscanctl -t /var/run/s6/services
        else
            printf "*** celery is not ready. Checking again in %2ds...(%d/%d)\n" $(((n**2+1)*5)) $n $m
        fi
    done
fi

echo "*** starting SJVA..."
s6-setuidgid abc \
    python sjva3.py \
        --repeat=${SJVA_REPEAT_COUNT} \
        --use_gevent=${USE_GEVENT} \
        --use_celery=${USE_CELERY} \
        --port=${SJVA_PORT}

EXITCODE=$?
echo "*** SJVA exited with code ${EXITCODE}"
if [ "$EXITCODE" = "0" ]; then
    echo '*** terminating all services...'
    s6-svscanctl -t /var/run/s6/services
else
    echo "*** restarting services..."
    printf $(expr $SJVA_REPEAT_COUNT + 1) > /var/run/s6/container_environment/SJVA_REPEAT_COUNT
    /etc/cont-init.d/20-config
    if [ "${USE_CELERY}" = "true" ]; then
        s6-svc -wr -r /var/run/s6/services/celery
    fi
fi
