# root/etc/services.d/celery/run
#!/usr/bin/with-contenv bash

if [ "${USE_CELERY}" = "true" ]; then
	echo "*** starting celery with ${CELERY_WORKER_COUNT} workers"	
	# # as root
	# export C_FORCE_ROOT='true'
	# exec celery -A sjva3.celery \
	# 	worker -c${CELERY_WORKER_COUNT} --loglevel=info

	# as root but with uid gid
	exec celery -A sjva3.celery \
		worker -c${CELERY_WORKER_COUNT} --loglevel=info --uid=${PUID:-911} --gid=${PGID:-911} > /dev/null 2>&1
fi

# dummy to demonize this service run
while true; do sleep 1h; done
