# root/etc/services.d/celery/run
#!/usr/bin/with-contenv bash

if [ "${USE_CELERY}" = "true" ]; then
	echo "*** starting celery with ${CELERY_WORKER_COUNT} workers"
	if [ ${PUID} -eq 0 ] && [ ${PGID} -eq 0 ]; then
		export C_FORCE_ROOT="true"
		export PYTHONWARNINGS="ignore::RuntimeWarning"
	fi
	exec s6-setuidgid abc \
		celery worker -c${CELERY_WORKER_COUNT} --loglevel=INFO
		#  > /dev/null 2>&1
fi
