# root/etc/services.d/redis/run
#!/usr/bin/with-contenv bash

if [ "${USE_CELERY}" = "true" ]; then
	echo "*** starting redis-server@${REDIS_PORT}"
	exec s6-setuidgid abc \
		redis-server --port ${REDIS_PORT} > /dev/null 2>&1
fi
