# root/etc/services.d/filebrowser/run
#!/usr/bin/with-contenv bash

if [ "${USE_FILEBROWSER}" = "true" ]; then
	echo "*** starting filebrowser"
	exec filebrowser -a 0.0.0.0 \
		-b "${FB_BASEURL:-/}" \
		-d "${FB_DATABASE:-/app/data/db/filebrowser.db}" \
		-p "${FB_PORT:-9998}" \
		-r "${FB_ROOT:-/}" \
		> /dev/null 2>&1
fi

# dummy to demonize this service run
while true; do sleep 1h; done
